# config/ai-dev-team.service
[Unit]
Description=AI Development Team MCP Server
After=network.target

[Service]
Type=simple
User=mark
WorkingDirectory=/home/mark/ai-development-team
EnvironmentFile=/home/mark/ai-development-team/.env
ExecStart=/home/mark/ai-development-team/start.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
